# files/example2/modules/vmxsquid/files/squid-softdog.cron
MAILTO=""
* * * * * root /usr/local/bin/squid-softdog
