# files/example2/modules/vmxsquid/files/vmx/sh-nj/run-puppetd.cron
MAILTO=""
10,30,50 *	* * *	root	/usr/bin/puppetd --debug --onetime
