# files/example2/modules/autodns/files/nagios-softdog.cron
MAILTO=""
* * * * * root /usr/local/bin/nagios-softdog

