# files/example2/modules/vmxcacti/files/count_user_rt_and_push.cron
MAILTO=""
*/2 * * * * root /usr/local/bin/count_user_rt_and_push.sh
