# files/example2/modules/vmxsquid/files/vmx/sh-zj/run-puppetd.cron
MAILTO=""
5,25,45 *	* * *	root	/usr/bin/puppetd --debug --onetime
