# files/example2/modules/autodns/files/internal/active_internal_autodns.cron
MAILTO=""
*/2 * * * * root /usr/local/bin/active_internal_autodns /var/log/nagios/status.dat
