# files/example2/modules/autodns/files/internal/vmx-pdnsd-softdog.cron
MAILTO=""
* * * * * root /usr/local/bin/vmx-pdnsd-softdog

