# files/example2/modules/autodns/files/external/auto_dns.cron
MAILTO=""
* * * * * root /usr/local/bin/active_auto_dns /var/log/nagios/status.dat|/usr/bin/nsupdate
