# files/example2/modules/autodns/files/external/start_stop_bind.cron
MAILTO=""
* * * * * root /usr/local/bin/start_stop_bind /var/log/nagios/status.dat|/bin/bash
