# files/example2/modules/vmxsquid/files/ido/sz-lg/run-puppetd.cron
MAILTO=""
15,35,55 *	* * *	root	/usr/bin/puppetd --debug --onetime
